# skbox-service/templates/skbox.service.j2
[Unit]
Description= Skillbox Service DevOps final

[Service]
Type=forking
EnvironmentFile=/etc/skillbox_env
ExecStart= /usr/bin/docker run {{ container_image }}

[Install]
WantedBy=multi-user.target